# argus.service
[Unit]
Description=Argus Ops Watchdog
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/perttu/argus

[Service]
Type=simple
User=perttu
Group=perttu
WorkingDirectory=~/argus
Environment="PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/etc/argus/argus.env
ExecStart=/usr/bin/bash -lc ./argus.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argus

# Resource limits — prevent runaway memory from LLM subprocess
MemoryMax=1G
MemoryHigh=768M

# Give the process time to finish a cycle on stop
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=%h/.openclaw/workspace/state %h/argus/logs

[Install]
WantedBy=multi-user.target
